# nginx/cron.d/certbot
0 0 * * * root /usr/bin/certbot renew >> /certbot.log
